# lib/tools/services/mongo-express/mongo-express.conf
description "Mongo-express server"
author      "Dave Kelsey" # I only wrote the upstart config

# used to be: start on startup
# until we found some mounts weren't ready yet while booting

start on started mountall
stop on shutdown

# automatically respawn

respawn
respawn limit 99 5

script

    export HOME="/home/vagrant/mongo-expres"
    exec su -s /bin/sh -c 'exec "$0" "$@"' vagrant -- /home/vagrant/mongo-express/bootstrap

end script

post-start script

   # optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"

end script
