# lib/tools/services/battlesnake/battlesnake.conf
description "BattleSnake server"
author      "SWU+ team"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting

start on started mountall
stop on shutdown

# automatically respawn

respawn
respawn limit 99 5

script

    export HOME="/home/vagrant/battlesnake"
    exec su -s /bin/sh -c 'exec "$0" "$@"' vagrant -- /home/vagrant/battlesnake/bin/bootstrap

end script

post-start script

   # optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "BattleSnake Server  has started!"

end script
